# assets/main.service
[Unit]
Description=

[Service]
WorkingDirectory=TMPDIR
ExecStart=TMPDIR/main
KillMode=process
Restart=on-failure
Type=exec

[Install]
WantedBy=multi-user.target